# deploy/yantresh-pull.service
# Install: sudo cp deploy/yantresh-pull.* /etc/systemd/system/
# Edit WorkingDirectory/ExecStart below if this repo isn't at /opt/yantresh-hub.
[Unit]
Description=Pull latest yantresh-hub images and recreate changed containers
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
# Runs as the user that owns /opt/yantresh-hub and the docker group, not
# root — DEPLOY.md step 5 rewrites this line to $USER before installing.
User=deploy
Group=docker
WorkingDirectory=/opt/yantresh-hub
ExecStart=/opt/yantresh-hub/deploy/pull-update.sh
NoNewPrivileges=true
